# build/supervisord-client.conf
[supervisord]
nodaemon=true

[program:proxy]
command=/app/proxy

[program:tcp]
command=/app/tcp

[program:udp]
command=/app/udp